# extras/midi-control.service
[Unit]
Description=Trigger actions on MIDI input.

[Service]
Type=simple
ExecStart=/usr/bin/midi-control

[Install]
WantedBy=default.target
